# sysrqd.service
[Unit]
Description=SysRqd Telnet Daemon
After=network.target

[Service]
Type=forking
Restart=always
ExecStart=-/usr/sbin/sysrqd
Nice=-19

[Install]
WantedBy=multi-user.target
